# mautic-builder/mautic.crontab
*/1 * * * * /usr/local/bin/php /var/www/html/app/console mautic:segments:update >> /var/log/cron.pipe 2>&1
*/1 * * * * /usr/local/bin/php /var/www/html/app/console mautic:campaigns:rebuild >> /var/log/cron.pipe 2>&1
*/1 * * * * /usr/local/bin/php /var/www/html/app/console mautic:campaigns:trigger >> /var/log/cron.pipe 2>&1